# services.d/zeroclaw/run
#!/command/with-contenv bash

set -e

echo "Starting zeroclaw..."

zeroclaw gateway